# supervisor.ini

[supervisord]
logfile = ./supervisor/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = False
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = ./supervisor/
nocleanup = true
childlogdir = ./supervisor/

[program:serve]
command=php artisan serve --verbose
process_name=%(program_name)s_%(process_num)02d
redirect_stderr=true
stdout_logfile=./supervisor/server-worker.log

[program:queue-worker]
command=php ./artisan queue:work --queue=default -vvv
numprocs=4
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=./supervisor/queue-worker.log

[program:queue-highprio-worker]
command=php ./artisan queue:work --queue=high -vvv
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=./supervisor/queue-highprio-worker.log

[program:reverb-worker]
command=php ./artisan reverb:start --debug
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
redirect_stderr=true
stdout_logfile=./supervisor/reverb-worker.log
